# services/miner-syncd/miner-syncd.service
[Unit]
Description=Miner-Lab Sync Daemon
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=minerlab
WorkingDirectory=/opt/miner-lab
EnvironmentFile=/etc/miner-lab/.env
ExecStart=/opt/miner-lab/venv/bin/python3 /opt/miner-lab/scripts/miner-syncd.py
Restart=on-failure
RestartSec=10
TimeoutStopSec=15
StandardOutput=append:/var/log/miner-lab/miner-syncd.log
StandardError=append:/var/log/miner-lab/miner-syncd.err

[Install]
WantedBy=multi-user.target